# srcs/wordpress/supervisord.conf
[supervisord]
nodaemon=true

[program:nginx]
command=nginx -g 'daemon off;'
autostart=true
autorestart=true
startsecs=0
redirect_stderr=true

[program:php-fpm7]
command=php-fpm7 -F -R
autostart=true
autorestart=true
startsecs=0
redirect_stderr=true
